# bin/old/gunicorn_rentality.service
[Unit]
Description=gunicorn daemon
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/home/Project/rentality
ExecStart=/home/Project/py3_env_rentality/bin/gunicorn --access-logfile - --workers 3 --bind unix:/home/Project/rentality/rentality_src.sock rentality.wsgi:application

[Install]
WantedBy=multi-user.target